# .config/systemd/user/plank.service
[Unit]
Description=Plank (Run Once)

[Service]
Type=oneshot
ExecStart=/usr/bin/plank

[Install]
WantedBy=default.target

